# little7-installer/systemd/lyra-telegram-bridge.service
[Unit]
Description=UnifAI Telegram Bridge (Command and Control)
After=network-online.target lyra-supervisor.service
Wants=network-online.target
Requires=lyra-supervisor.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/little7/supervisor

ExecStart=/usr/bin/python3 /opt/little7/supervisor/plugins/telegram_bridge/bot_listener.py

# Authorized operator and optional runtime overrides
EnvironmentFile=-/etc/little7/telegram_bridge.env

Restart=always
RestartSec=5
ProtectSystem=full
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lyra-telegram-bridge

[Install]
WantedBy=multi-user.target
WantedBy=multi-user.target